# packaging/app-network-report.cron
*/5 * * * * root /usr/sbin/network2db >/dev/null 2>&1
12 3 * * * root /usr/sbin/purge-report-tables network >/dev/null 2>&1
